# tendwired.service.example
[Unit]
Description=Tendwire daemon

[Service]
Type=simple
UMask=0077
Environment=TENDWIRE_HERDR_BACKEND=socket
Environment=TENDWIRE_DB_PATH=%h/.local/share/tendwire/tendwire.db
Environment=TENDWIRE_TURN_REFRESH_INTERVAL_SECONDS=2.0
Environment=TENDWIRE_TURN_REFRESH_WORKERS=4
ExecStart=%h/.local/bin/tendwire daemon --db-path %h/.local/share/tendwire/tendwire.db
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=20s
FinalKillSignal=SIGKILL
Restart=always
RestartSec=5s

[Install]
WantedBy=default.target
